# crontab
0 3 * * 1 root /./borg/backup.sh >> /var/log/cron.log 2>&1
